# .config/autostart/ssh-add.desktop
[Desktop Entry]
Exec=ssh-add -q [ssh keys]
Name=ssh-add
Type=Application