# .config/autostart/ssh-add.desktop
[Desktop Entry]
Exec=ssh-add -q /home/canalejas/.ssh/id_rsa
Name=ssh-add
Type=Application
